# vultisig-agent.service.tmpl
# Systemd service template for Vultisig Dev Agent.
# Generated by setup.sh — or copy and edit manually.
#
# WHY SYSTEMD: Without it, if the VPS reboots (kernel updates, provider
# maintenance) or the process crashes (OOM, unexpected error), the agent
# dies silently. No issues get picked up until someone manually SSHs in.
# Systemd auto-starts on boot, auto-restarts on crash (after 10s cooldown),
# and collects logs via journalctl.
#
# Install:
#   sudo cp vultisig-agent.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now vultisig-agent
#
# Logs:
#   sudo journalctl -u vultisig-agent -f        # live tail
#   sudo journalctl -u vultisig-agent --since today  # today's logs

[Unit]
Description=Vultisig Dev Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=REPLACE_USER
WorkingDirectory=REPLACE_WORKING_DIR
ExecStart=REPLACE_OPENCLAW_PATH start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vultisig-agent

[Install]
WantedBy=multi-user.target
